# gencore.socket
[Unit]
Description=Gencore socket

[Socket]
ListenStream=/var/run/gencored.socket
Accept=true

[Install]
WantedBy=socket.target

